# nodepool-scripts/static-network-config
# Copyright (c) 2015 Hewlett-Packard Development Company, L.P.
# chkconfig: 2345 05 95
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#    http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or
# implied.
#
# See the License for the specific language governing permissions and
# limitations under the License.

if [ ! -z $1 -a $1 != start ] ; then
    exit 0
fi

# Start by putting down base scripts
if [ -d /etc/sysconfig/network-scripts ]; then
    # Bail and don't do anything if we have one - handles reboot
    if [ -f /etc/sysconfig/network-scripts/ifcfg-eth0 ] ; then
        exit 0
    fi
    for interface in eth0 eth1 ; do
        cat >/etc/sysconfig/network-scripts/ifcfg-$interface <<EOF
DEVICE=$interface
BOOTPROTO=dhcp
ONBOOT=on
EOF
    done
elif [ -d /etc/network ]; then
    # Bail and don't do anything if we have one - handles reboot
    if grep eth0 /etc/network/interfaces >/dev/null 2>&1 ; then
        exit 0
    fi
    for interface in eth0 eth1 ; do
        cat >>/etc/network/interfaces <<EOF
auto $interface
iface $interface inet dhcp
EOF
    done
fi

# Test to see if config-drive exists. If not, skip and assume DHCP networking
# will work becasue sanity
if ! blkid -t LABEL="config-2" ; then
    exit 0
fi

# Mount config drive- but check to make sure it's Rackspace and not config
# drive from somewhere else
mkdir -p /mnt/config
mount /dev/disk/by-label/config-2 /mnt/config
grep -i rackspace /mnt/config/openstack/latest/vendor_data.json >/dev/null 2>&1 || exit 0

python /usr/local/bin/read-vendor-json.py

if [ -f /etc/network/interfaces ] ; then
    ifdown -a ; ifup -a
else
    for interface in /etc/sysconfig/network-scripts/ifcfg-eth* ; do
        name=$(echo $interface | cut -f3 -d-)
        ifdown $name ; ifup $name
    done
fi
